# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:term-app]
command=node /app/index.js
autorestart=true
environment=NODE_ENV=production,PORT=8999
