# control-node/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:webshell]
command=node /app/build/main.js
autorestart=true
user=root
